# dns/dns/service/dnsmasq/run
#!/bin/sh
exec /usr/sbin/dnsmasq --keep-in-foreground --user=root --local-service --cache-size=2000 --no-negcache --local-ttl=3600 >&2
